# etc/crontabs/root
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
HOME=/

# run-parts
* */2 * * * echo > /var/log/cron.log
