# etc/inittab
::sysinit:/bin/cttyhack /etc/rc.d/start
ttyS1::respawn:/bin/cttyhack /bin/sh
