# etc/inittab
::sysinit:/bin/cttyhack /etc/init.d/rcS
